# smf/manifests/ufds-capi-8082.xml.in
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="smartdc-ufds-capi">
    <service name="smartdc/application/ufds-capi-8082" type="service" version="1">
        <dependency name="network"
                    grouping="require_all"
                    restart_on="error"
                    type="service">
	    <service_fmri value="svc:/network/physical" />
        </dependency>

        <dependency name="filesystem"
                    grouping="require_all"
                    restart_on="error"
                    type="service">
	    <service_fmri value="svc:/system/filesystem/local" />
        </dependency>

        <dependency name="ufds"
                    grouping="require_all"
                    restart_on="error"
                    type="service">
	    <service_fmri value="svc:/smartdc/application/ufds-master" />
        </dependency>

        <exec_method type="method"
                     name="start"
                     exec="/opt/smartdc/ufds/build/node/bin/node capi.js -f etc/config.json -d 2 &amp;"
                     timeout_seconds="30">
            <method_context working_directory="/opt/smartdc/ufds">
                <method_credential user="nobody"
                                   group="nobody"
                                   privileges="basic,net_privaddr"/>
                <method_environment>
                    <envvar name="PATH"
                            value="/opt/smartdc/ufds/build/node/bin:/opt/local/bin:/usr/bin:/usr/sbin:/bin"/>
                    <envvar name="LD_PRELOAD_32"
                            value="/usr/lib/extendedFILE.so.1" />
                    <envvar name="PORT"
                            value="8082" />
                </method_environment>
            </method_context>
        </exec_method>

        <exec_method type="method"
                     name="stop"
                     exec=":kill"
                     timeout_seconds="30" />

        <property_group name="application" type="application" />
        <instance name="default" enabled="true" />
        <stability value='Unstable' />

    <template>
	<common_name><loctext xml:lang="C">UFDS Master Service</loctext></common_name>
    </template>

    </service>
</service_bundle>
